# systemd/sort_images.timer
[Unit]
Description=Sort Images

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
