# systemd/sort_images.service
[Unit]
Description=SortImages Service

[Service]
Type=oneshot
ExecStart=/home/fargus/Projects/SortImages/main.py
